# deploy/systemd/xgg-backup.timer
[Unit]
Description=XGenGuardian Postgres backup — daily 03:00 UTC
Documentation=file:///home/xgenguardian/code/docs/runbooks/postgres-backup.md

[Timer]
# Fire every day at 03:00 UTC
OnCalendar=*-*-* 03:00:00 UTC
# Catch up if the system was off when the timer should have fired
Persistent=true
# Small random delay to avoid thundering-herd on multi-node setups
RandomizedDelaySec=300
Unit=xgg-backup.service

[Install]
WantedBy=timers.target
